# .crontab
 0 14 * * * brew update && brew upgrade && brew cask upgrade
 0 14 * * * "$HOME"/config/tools/update.sh
 0 14 * * * rm -f "$HOME"/Desktop/Screen\ Shot* "$HOME"/Downloads/*.gif "$HOME"/Downloads/*.mkv "$HOME"/Downloads/*.mp4 "$HOME"/Downloads/*.dmg

